# docker/lnmp/schedule.cron
* * * * * root bash -l -c "php /opt/ci123/www/html/artisan schedule:run"
